# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable trident.service
systemctl restart trident.service
