# packaging/scripts/preremove.sh
systemctl stop trident.service
systemctl disable trident.service
